# player-1/supervisor/apache2.conf
[program:apache2]
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true